# quadlet/launcherd.container
# launcherd.container — launch-controller daemon (the --launcher and --dispatch doors)
#
# UNLIKE keeperd/netd/scoutd, launcherd is NOT ready to run from this unit
# yet — Image=localhost/launcherd:dev does not exist as a build target (no
# `launcherd-image` derivation in flake.nix, unlike keeperd-image/netd-image/
# scoutd-image). Building one is a real prerequisite, not just "cp this file
# and go": it needs the `podman` CLI binary inside the image (launcherd
# shells out to `podman run`/`podman inspect`/`podman kill` directly — see
# launcherd.ts's buildPodmanArgv/getContainerId), plus a way to reach the
# host's container runtime from inside its own container (see the podman
# socket note below). Treat this unit as the target shape to build toward,
# not a drop-in-and-run file.
#
# Boxes ask to spawn/inspect/kill sub-boxes through /run/doors/launcherd.sock
# (the --launcher door, attenuated, name-based identity — see LAUNCHERD.md),
# or request an independent task session through /run/doors/dispatch.sock
# (the --dispatch door, doors-blind, allow-listed rooms only — see
# launcherd.ts's DISPATCH_METHODS / handleDispatch). launcherd owns the
# runtime for BOTH lanes; boxes never hold podman.
#
# Install:
#   cp *.volume *.container ~/.config/containers/systemd/
#   systemctl --user daemon-reload
#   systemctl --user enable --now launcherd
#
# Logs:
#   journalctl --user -u launcherd -f

[Unit]
Description=launcherd — launch-controller daemon for claude-box (launch + dispatch)
After=claude-doors-volume.service

[Container]
Image=localhost/launcherd:dev
ContainerName=launcherd

# Socket directory — same shared volume every other door uses. launcherd
# listens on BOTH launcherd.sock (general) and dispatch.sock (narrow) here —
# see launcherd.ts's serve(socketPath, dispatchSocketPath).
Volume=%h/.claude-box/run:/run/doors

# launcherd is the one daemon that must reach the HOST's container runtime
# (to actually spawn/inspect/kill sibling containers) — unlike keeperd/netd/
# scoutd, which are pure socket daemons with no host-control surface at all.
# RECOMMENDED approach, not yet validated against a real deployment: mount
# the rootless podman API socket and point the in-image `podman` CLI at it
# over the wire, rather than trying to share a container runtime directly.
#   Volume=%t/podman/podman.sock:/run/podman/podman.sock
#   Environment=CONTAINER_HOST=unix:///run/podman/podman.sock
# This is a meaningfully more privileged grant than anything keeperd/netd/
# scoutd need (it's adjacent to "control the host's container runtime") —
# treat the exact socket path, permissions, and whether a narrower proxy is
# warranted as an open question for security review before enabling this in
# a real deployment, not a settled default.
Volume=%t/podman/podman.sock:/run/podman/podman.sock
Environment=CONTAINER_HOST=unix:///run/podman/podman.sock

# Security hardening — same floor as the other doors. NoNewPrivileges/
# DropCapability=all/ReadOnly still apply to launcherd's OWN process; they do
# NOT reduce what the mounted podman socket itself can do once reached (see
# the note above).
NoNewPrivileges=true
ReadOnly=true
ReadOnlyTmpfs=true
DropCapability=all
SecurityLabelDisable=true

# Resource limits
PidsLimit=256
Memory=512m

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
