# quadlet/netd.container
# netd.container — allowlist egress proxy (the --net door)
#
# Boxes route HTTPS_PROXY through /run/doors/netd.sock.
# netd enforces the destination allowlist; boxes have no direct egress.
#
# Install:
#   systemctl --user enable --now netd
#
# Logs:
#   journalctl --user -u netd -f

[Unit]
Description=netd — allowlist egress proxy for claude-box
After=claude-doors-volume.service

[Container]
Image=localhost/netd:dev
ContainerName=netd

# Socket directory — mount host's ~/.claude-box/run so sockets are accessible
# from both the host CLI and containers. Adjust path for your home directory.
Volume=/Users/bobby/.claude-box/run:/run/doors

# netd DOES need network (it's the egress door)
# Default: uses podman's default network
# Uncomment for explicit control:
# Network=slirp4netns

# Security hardening
NoNewPrivileges=true
ReadOnly=true
ReadOnlyTmpfs=true
DropCapability=all
SecurityLabelDisable=true

# Resource limits
PidsLimit=128
Memory=256m

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
