# quadlet/remote-serve.container
# remote-serve.container — the singleton RC dispatcher bastion
#
# UNLIKE `claude-box --remote-serve`'s CLI path, this systemd-started
# container has no invocation moment for claude-box.ts to mint a fresh
# signed authd grant or build the RC bootstrap script inline (mintAuthGrant/
# buildRemoteServeScript in run()). Both are produced host-side, once per
# start, by ExecStartPre=, and handed to the container two ways:
#   - the grant: base64, via an EnvironmentFile= (see cmdMintAuthGrant)
#   - the boot script: a plain file, bind-mounted read-only (see
#     cmdPrintRcBootScript) — NOT embedded directly in this unit's Exec=,
#     because systemd's own argv-splitting rules (distinct from POSIX
#     shell) make a script this size and this quote-heavy genuinely risky
#     to hand-transcribe correctly without a real host to test against.
#     Referencing it by path keeps Exec= itself trivial.
# The signing key never leaves the host either way (see lib/box-keys.ts) —
# only the signed grant object does, same as the CLI path.
#
# This unit is unvalidated against a real systemd+podman host (developed and
# tested only in a sandboxed environment with no such host available) —
# the mechanism and every piece it depends on (cmdMintAuthGrant,
# cmdPrintRcBootScript, authLeaseFromEnvCmd) are unit-tested in isolation,
# but the ExecStartPre=/EnvironmentFile=/bind-mounted-script wiring here has
# not been exercised end-to-end. Sanity-check `systemctl --user start
# remote-serve` + `journalctl --user -u remote-serve` on a real deployment
# before relying on it.
#
# Install:
#   cp *.volume *.container ~/.config/containers/systemd/
#   systemctl --user daemon-reload
#   systemctl --user enable --now remote-serve
#
# Logs:
#   journalctl --user -u remote-serve -f
#
# Attaching locally (the Claude app is the normal way in — this is only for
# local debugging): `podman attach claude-box-remote-serve`. No `-it` here
# (see below), so there's no local tty by default the way the manual
# `claude-box --remote-serve` CLI entrypoint still has.

[Unit]
Description=remote-serve — the singleton claude-box RC dispatcher bastion
Requires=claude-doors-volume.service launcherd.service
After=claude-doors-volume.service launcherd.service

[Container]
Image=localhost/claude-personal:dev
ContainerName=claude-box-remote-serve

# Same shared socket dir as every door — this box holds net+auth+dispatch
# (see claude-box.ts's --remote-serve block in planLaunch).
Volume=%h/.claude-box/run:/run/doors

# The freshly-minted grant (see [Service] ExecStartPre= below), read by the
# boot script's lease step (authLeaseFromEnvCmd) at $CLAUDE_BOX_RC_GRANT.
EnvironmentFile=%t/claude-box-rc-grant.env

# The boot script itself (see [Service] ExecStartPre= below) — generated
# fresh each start (so it always matches the current claude-box code), not
# baked into the image.
Volume=%t/claude-box-rc-boot.sh:/rc-boot.sh:ro

# RC needs egress (leasing the credential from authd, then ordinary
# inference/tool traffic) — default podman network, same posture as netd.
# No Network=none here.

# The credential is LEASED from authd, never a persisted refresh token (see
# claude-box.ts's run() comment on this exact point) — a throwaway tmpfs
# config dir, never the shared claude-config volume, so nothing about this
# bastion's identity lives at rest. Mirrors run()'s own
# `--tmpfs BOX_CONFIG_DIR:rw,mode=1777`.
Tmpfs=/home/claude/.config/claude:rw,mode=1777

# Boot straight into the generated script — see the ExecStartPre= note above
# for why this is a bind-mounted file rather than inline script content.
Entrypoint=sh
Exec=/rc-boot.sh claude-box remote-control --name dispatch --remote-control-session-name-prefix claude-box --spawn session

# Deliberately NO "-it" here: the Claude app attaches over claude.ai/code,
# never a local tty, so a Quadlet-managed bastion runs detached like every
# other unit in this directory. The manual `claude-box --remote-serve` CLI
# entrypoint is unaffected and keeps -it for local interactive use — the
# same dev-vs-prod split keeperd/scoutd already have between `doors serve`
# and their own .container units.

# Security hardening
NoNewPrivileges=true
DropCapability=all
SecurityLabelDisable=true

# Resource limits
PidsLimit=2048
Memory=1g

[Service]
# Mint this start's grant + boot script host-side (see the [Unit] comment
# above for why). %t is the same systemd runtime-dir specifier used by the
# Volume=/EnvironmentFile= lines above, so both sides agree on the path
# without any per-machine editing.
ExecStartPre=/bin/sh -c 'claude-box internal-mint-auth-grant --audience claude-box-remote-serve | (echo -n "CLAUDE_BOX_RC_GRANT=" && cat) > %t/claude-box-rc-grant.env'
ExecStartPre=/bin/sh -c 'claude-box internal-print-rc-boot-script > %t/claude-box-rc-boot.sh && chmod +r %t/claude-box-rc-boot.sh'

Restart=always
RestartSec=5

[Install]
WantedBy=default.target
