# quadlet/scoutd.container
# scoutd.container — external read daemon (the --scout door)
#
# Boxes request repo/PR/URL content through /run/doors/scoutd.sock.
# scoutd holds read tokens; boxes hold nothing and have no network.
#
# Install:
#   systemctl --user enable --now scoutd
#
# Logs:
#   journalctl --user -u scoutd -f

[Unit]
Description=scoutd — external read daemon for claude-box
After=claude-doors-volume.service

[Container]
Image=localhost/scoutd:dev
ContainerName=scoutd

# Socket directory — mount host's ~/.claude-box/run so sockets are accessible
# from both the host CLI and containers. Adjust path for your home directory.
Volume=/Users/bobby/.claude-box/run:/run/doors:U

# Token/credential volume (scoutd-only)
# Volume=claude-scout-creds.volume:/creds:ro

# scoutd needs network (it fetches external content)
# Default: uses podman's default network
# Uncomment for explicit control:
# Network=slirp4netns

# Security hardening
NoNewPrivileges=true
ReadOnly=true
ReadOnlyTmpfs=true
DropCapability=all
SecurityLabelDisable=true

# Resource limits
PidsLimit=128
Memory=512m

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
