# quadlet/scoutd.container
# scoutd.container — external read daemon (the --scout door)
#
# Boxes request repo/PR/URL content through /run/doors/scoutd.sock.
# scoutd holds read tokens; boxes hold nothing and have no network.
#
# Install:
#   systemctl --user enable --now scoutd
#
# Logs:
#   journalctl --user -u scoutd -f

[Unit]
Description=scoutd — external read daemon for claude-box
# scoutd egresses through scout-netd (its entrypoint relays loopback→
# scout-netd.sock), so that door must be up first.
After=claude-doors-volume.service scout-netd.service
Requires=scout-netd.service

[Container]
Image=localhost/scoutd:dev
ContainerName=scoutd

# Socket directory — mount the host's ~/.claude-box/run so sockets are accessible
# from both the host CLI and containers. %h expands to the running user's home
# (systemd specifier), so this stays host-agnostic — no per-machine editing.
Volume=%h/.claude-box/run:/run/doors

# Token/credential volume (scoutd-only)
# Volume=claude-scout-creds.volume:/creds:ro

# scoutd holds NO NIC. Its egress is forced through the scout-netd door: the
# image entrypoint bridges loopback→/run/doors/scout-netd.sock and sets
# SCOUTD_PROXY. netd is the egress chokepoint — interposition, not cooperation.
Network=none

# Security hardening
NoNewPrivileges=true
ReadOnly=true
ReadOnlyTmpfs=true
DropCapability=all
SecurityLabelDisable=true

# Resource limits
PidsLimit=128
Memory=512m

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
